# scripts/rpm/postinstall.in
#
# Copyright (C) 2022 The Falco Authors.
#
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#

chosen_driver=
enable_falcoctl=1

if [ $1 -eq 1 ]; then
        if [ -x /usr/bin/dialog ]; then
            # If dialog is installed, create a dialog to let users choose the correct driver for them
            CHOICE=$(dialog --clear --title "Falco drivers" --menu "Choose your preferred driver:" 12 40 5 \
                    1 "Don't start" \
                    2 "Kmod" \
                    3 "eBPF" \
                    4 "Modern eBPF" \
                    5 "Plugin" \
                    2>&1 >/dev/tty)
                case $CHOICE in
                    2)
                        chosen_driver="kmod"
                        ;;
                    3)
                        chosen_driver="bpf"
                        ;;
                    4)
                        chosen_driver="modern-bpf"
                        ;;
                    5)
                        chosen_driver="plugin (needs manual configuration)"
                        ;;
                esac
                CHOICE=$(dialog --clear --title "Falcoctl" --menu "Do you want to follow automatic ruleset updates?" 10 40 2 \
                        1 "Yes" \
                        2 "No" \
                        2>&1 >/dev/tty)
                case $CHOICE in
                    2)
                    enable_falcoctl=0
                    ;;
                esac
                clear
        fi
fi

set -e

# If needed, try to load/compile the driver through falco-driver-loader
case "$chosen_driver" in
    "kmod")
      echo "[POST-INSTALL] Call 'falco-driver-loader module':"
      falco-driver-loader module
      ;;
    "bpf")
      echo "[POST-INSTALL] Call 'falco-driver-loader bpf':"
      falco-driver-loader bpf
      ;;
esac

# validate rpm macros by `rpm -qp --scripts <rpm>`
# RPM scriptlets: https://docs.fedoraproject.org/en-US/packaging-guidelines/Scriptlets/#_systemd
#                 https://docs.fedoraproject.org/en-US/packaging-guidelines/Scriptlets/#_syntax

# systemd_post macro expands to
# if postinst:
#   `systemd-update-helper install-system-units <service>`
%systemd_post "falco-$chosen_driver.service"

# post install mirrored from .deb
if [ $1 -eq 1 ]; then
        if [ -n "$chosen_driver" ]; then
            echo "[POST-INSTALL] Enable 'falco-$chosen_driver.service':"
            systemctl --system enable "falco-$chosen_driver.service" || true
            echo "[POST-INSTALL] Start 'falco-$chosen_driver.service':"
            systemctl --system start "falco-$chosen_driver.service" || true
        fi
        if [ $enable_falcoctl -eq 1 ]; then
          echo "[POST-INSTALL] Enable 'falcoctl-artifact-follow.service':"
          systemctl --system enable "falcoctl-artifact-follow.service" || true
          echo "[POST-INSTALL] Start 'falcoctl-artifact-follow.service':"
          systemctl --system start "falcoctl-artifact-follow.service" || true
        fi
fi

# post upgrade mirrored from .deb
if [ $1 -gt 1 ]; then
    if [ -d /run/systemd/system ]; then
        echo "[POST-INSTALL] Trigger deamon-reload:"
        systemctl --system daemon-reload || true
        if [ -n "$chosen_driver" ]; then
          echo "[POST-INSTALL] Trigger 'falco-$chosen_driver.service' condrestart:"
          # restart falco on upgrade if service is already running
          systemctl --system condrestart "falco-$chosen_driver.service" || true
        fi
        if [ $enable_falcoctl -eq 1 ]; then
          echo "[POST-INSTALL] Trigger 'falcoctl-artifact-follow.service' condrestart:"
          # restart falcoctl-artifact-follow on upgrade if service is already running
          systemctl --system condrestart "falcoctl-artifact-follow.service" || true
        fi
    fi
fi
