# scripts/rpm/postinstall.in
#
# Copyright (C) 2022 The Falco Authors.
#
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
set -e

chosen_driver=

case "$1" in
    configure)
        if [ -x /usr/bin/dialog ]; then
            # If dialog is installed, create a dialog to let users choose the correct driver for them
            CHOICE=$(dialog --clear --backtitle "Choose your preferred driver" --title "Falco driver" --menu "Choose one of the following options:" 15 40 4 \
                    1 "Don't start" \
                    2 "Kmod" \
                    3 "eBPF" \
                    4 "Modern eBPF" \
                    5 "Plugin" \
                    2>&1 >/dev/tty)
                clear
                case $CHOICE in
                    2)
                        chosen_driver="kmod"
                        ;;
                    3)
                        chosen_driver="bpf"
                        ;;
                    4)
                        chosen_driver="modern-bpf"
                        ;;
                    5)
                        chosen_driver="plugin"
                        ;;
                esac
        fi
    ;;
esac

# If needed, try to load/compile the driver through falco-driver-loader
case "$chosen_driver" in
    "kmod")
      falco-driver-loader module
      ;;
    "bpf")
      falco-driver-loader bpf
      ;;
esac

# validate rpm macros by `rpm -qp --scripts <rpm>`
# RPM scriptlets: https://docs.fedoraproject.org/en-US/packaging-guidelines/Scriptlets/#_systemd
#                 https://docs.fedoraproject.org/en-US/packaging-guidelines/Scriptlets/#_syntax

# systemd_post macro expands to
# if postinst:
#   `systemd-update-helper install-system-units <service>`
%systemd_post "falco@$chosen_driver.target"

# post install mirrored from .deb
if [ $1 -eq 1 ]; then
        # This will only remove masks created on package removal.
        /usr/bin/systemctl --system unmask "falco@.target" >/dev/null || true
        /usr/bin/systemctl --system unmask "falco-kmod.service" >/dev/null || true
        /usr/bin/systemctl --system unmask "falco-bpf.service" >/dev/null || true
        /usr/bin/systemctl --system unmask "falco-modern-bpf.service" >/dev/null || true
        /usr/bin/systemctl --system unmask "falco-plugin.service" >/dev/null || true

        if [ -n "$chosen_driver" ]; then
            # enable falco on installation
            # note: DEB postinstall script checks for changed symlinks
            /usr/bin/systemctl --system enable "falco@$chosen_driver.target" >/dev/null || true

            # start falco on installation
            /usr/bin/systemctl --system start "falco@$chosen_driver.target" >/dev/null || true
        fi
fi

# post upgrade mirrored from .deb
if [ $1 -gt 1 ]; then
    if [ -d /run/systemd/system ]; then
        /usr/bin/systemctl --system daemon-reload >/dev/null || true
        if [ -n "$chosen_driver" ]; then
            # restart falco on upgrade if service is already running
            /usr/bin/systemctl --system condrestart "falco@$chosen_driver.target" >/dev/null || true
        fi
    fi
fi
