# scripts/systemd/falco.target
[Unit]
Description=Main Falco target
Requires=multi-user.target
After=multi-user.target

[Install]
WantedBy=multi-user.target