# scripts/debian/postinst.in
#!/bin/sh
#
# Copyright (C) 2022 The Falco Authors.
#
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
set -e

chosen_driver=

case "$1" in
    configure)
        if [ -x /usr/bin/dialog ]; then
          # If dialog is installed, create a dialog to let users choose the correct driver for them
          CHOICE=$(dialog --clear --backtitle "Choose your preferred driver" --title "Falco driver" --menu "Choose one of the following options:" 15 40 4 \
                1 "Don't start" \
                2 "Kmod" \
                3 "eBPF" \
                4 "Plugin" \
                2>&1 >/dev/tty)
            clear
            case $CHOICE in
                1)
                    chosen_driver="kmod"
                    ;;
                2)
                    chosen_driver="ebpf"
                    ;;
                3)
                    chosen_driver="plugin"
                    ;;
            esac
        else
          # Default at old behavior: enable kmod driver
          chosen_driver="kmod"
        fi
    ;;
esac

# If needed, try to load/compile the driver through falco-driver-loader
case "$chosen_driver" in
    "kmod")
      falco-driver-loader module
      ;;
    "ebpf")
      falco-driver-loader bpf
      ;;
esac

# Based off what debhelper dh_systemd_enable/13.3.4 would have added
# ref: https://www.debian.org/doc/manuals/debmake-doc/ch05.en.html#debhelper

if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
        # This will only remove masks created by d-s-h on package removal.
        deb-systemd-helper unmask "falco@kmod.target" >/dev/null || true
        deb-systemd-helper unmask "falco@ebpf.target" >/dev/null || true
        deb-systemd-helper unmask "falco@plugin.target" >/dev/null || true

        # was-enabled defaults to true, so new installations run enable.
        if [ -n "$chosen_driver" ]; then
            if deb-systemd-helper --quiet was-enabled "falco@$chosen_driver.target"; then
                    # Enables the unit on first installation, creates new
                    # symlinks on upgrades if the unit file has changed.
                    deb-systemd-helper enable "falco@$chosen_driver.target" >/dev/null || true
            else
                    # Update the statefile to add new symlinks (if any), which need to be
                    # cleaned up on purge. Also remove old symlinks.
                    deb-systemd-helper update-state "falco@$chosen_driver.target" >/dev/null || true
            fi
        fi
fi

if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
    if [ -d /run/systemd/system ]; then
      systemctl --system daemon-reload >/dev/null || true
      if [ -n "$chosen_driver" ]; then
        if [ -n "$2" ]; then
          _dh_action=restart
        else
          _dh_action=start
        fi
        deb-systemd-invoke $_dh_action "falco@$chosen_driver.target" >/dev/null || true
      fi
    fi
fi
