# scripts/debian/postrm.in
#!/bin/sh
#
# Copyright (C) 2022 The Falco Authors.
#
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#

# Based off what debhelper dh_systemd_enable/13.3.4 would have added
# ref: https://www.debian.org/doc/manuals/debmake-doc/ch05.en.html#debhelper

set -e

if [ -d /run/systemd/system ] && [ "$1" = remove ]; then
        systemctl --system daemon-reload >/dev/null || true
fi

if [ "$1" = "remove" ]; then
        if [ -x "/usr/bin/deb-systemd-helper" ]; then
                systemctl --system mask 'falco@.target' >/dev/null || true
                systemctl --system mask 'falco-kmod.service' >/dev/null || true
                systemctl --system mask 'falco-bpf.service' >/dev/null || true
                systemctl --system mask 'falco-modern-bpf.service' >/dev/null || true
                systemctl --system mask "falco-plugin.service" >/dev/null || true
        fi
fi