# cmake-build-debug/scripts/debian/prerm
#!/bin/sh
#
# Copyright (C) 2019 The Falco Authors.
#
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
set -e

NAME="falco"

if [ -x "/etc/init.d/$NAME" ]; then
		if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
				invoke-rc.d $NAME stop || exit $?
		else
				/etc/init.d/$NAME stop || exit $?
		fi
fi

DKMS_PACKAGE_NAME="falco"
DKMS_VERSION="96bd9bc560f67742738eb7255aeb4d03046b8045"

case "$1" in
	remove|upgrade|deconfigure)
		if [  "$(dkms status -m $DKMS_PACKAGE_NAME -v $DKMS_VERSION)" ]; then
			dkms remove -m $DKMS_PACKAGE_NAME -v $DKMS_VERSION --all
		fi
	;;
esac
